# file.list
etc/cron.d/vpsnetclean
etc/cron.d/vpsreboot
etc/ifcfg-venet0
etc/ifdown-venet
etc/ifup-venet
etc/ve-light.conf-sample
etc/ve-light.conf-sample.ia64
etc/ve-light.conf-sample.x86_64
etc/ve-vps.basic.conf-sample
etc/ve-vps.basic.conf-sample.ia64
etc/ve-vps.basic.conf-sample.x86_64
etc/vpsnetclean
etc/vpsreboot
etc/vps-create
etc/vps-functions
etc/vps-net_add
etc/vps-net_del
etc/vps-postcreate
etc/vps-stop
etc/vz
etc/dists/scripts/debian-add_ip.sh
etc/dists/scripts/debian-del_ip.sh
etc/dists/scripts/debian-set_hostname.sh
etc/dists/scripts/functions
etc/dists/scripts/gentoo-add_ip.sh
etc/dists/scripts/gentoo-del_ip.sh
etc/dists/scripts/gentoo-set_hostname.sh
etc/dists/scripts/gentoo-set_ugid_quota.sh
etc/dists/scripts/postcreate.sh
etc/dists/scripts/redhat-7.x_add_ip.sh
etc/dists/scripts/redhat-add_ip.sh
etc/dists/scripts/redhat-del_ip.sh
etc/dists/scripts/redhat-set_hostname.sh
etc/dists/scripts/set_dns.sh
etc/dists/scripts/set_ugid_quota.sh
etc/dists/scripts/set_userpass.sh
etc/dists/scripts/slackware-9.x_add_ip.sh
etc/dists/scripts/slackware-9.x_del_ip.sh
etc/dists/scripts/slackware-add_ip.sh
etc/dists/scripts/slackware-del_ip.sh
etc/dists/scripts/slackware-set_hostname.sh
etc/dists/scripts/slackware-set_ugid_quota.sh
etc/dists/scripts/suse-7.x_add_ip.sh
etc/dists/scripts/suse-8.x_add_ip.sh
etc/dists/scripts/suse-8.x_del_ip.sh
etc/dists/scripts/suse-add_ip.sh
etc/dists/scripts/suse-del_ip.sh
etc/dists/scripts/suse-set_hostname.sh
etc/dists/centos.conf
etc/dists/debian.conf
etc/dists/distribution.conf-template
etc/dists/fedora-core.conf
etc/dists/gentoo.conf
etc/dists/mandrake.conf
etc/dists/redhat-7.3.conf
etc/dists/redhat.conf
etc/dists/slackware-9.0.conf
etc/dists/slackware.conf
etc/dists/sles.conf
etc/dists/suse-7.3.conf
etc/dists/suse-8.0.conf
etc/dists/suse.conf
etc/dists/opensuse.conf
etc/dists/rhel-3.conf
etc/dists/default
etc/dists/Makefile
etc/init.d/vz
etc/logrotate.d/vzctl
man/Makefile
man/arpsend.8
man/vps.conf.5
man/vz.5
man/vzcalc.8
man/vzcfgvalidate.8
man/vzcpucheck.8
man/vzctl.8
man/vzlist.8
man/vzmemcheck.8
man/vzpid.8
man/vzsplit.8
Makefile
vzcpucheck
vzctl.spec
vzpid
include/config.h
include/cap.h
include/cpu.h
include/create.h
include/dev.h
include/dist.h
include/env.h
include/exec.h
include/fs.h
include/iptables.h
include/list.h
include/lock.h
include/logger.h
include/modules.h
include/net.h
include/quota.h
include/readelf.h
include/res.h
include/script.h
include/types.h
include/ub.h
include/util.h
include/validate.h
include/vzctl.h
include/vzctl_param.h
include/vzerror.h
include/vzlist.h
src/lib/config.c
src/lib/cap.c
src/lib/create.c
src/lib/fs_bind.c
src/lib/fs_simfs.c
src/lib/cpu.c
src/lib/dev.c
src/lib/dist.c
src/lib/env.c
src/lib/exec.c
src/lib/fs.c
src/lib/iptables.c
src/lib/list.c
src/lib/lock.c
src/lib/logger.c
src/lib/modules.c
src/lib/net.c
src/lib/quota.c
src/lib/readelf.c
src/lib/res.c
src/lib/script.c
src/lib/ub.c
src/lib/util.c
src/Makefile
src/arpsend.c
src/enter.c
src/modules.c
src/validate.c
src/vzcalc.c
src/vzcfgvalidate.c
src/vzctl-actions.c
src/vzctl.c
src/vzlist.c
src/vzmemcheck.c
src/vzsplit.c
LICENSE.QPL
ChangeLog
