# etc/cron.d/vz.in
# Cron jobs needed for OpenVZ

# -- WARNING -- WARNING -- WARNING -- WARNING -- WARNING -- WARNING --
# ------------------  EDIT THE CORRECT FILE  -------------------------
# This file is copied to @DST_CRONDDIR@/vz from
# @CRONDDIR@/vz when the vz service is started
# via its init.d script and the file
# @DST_CRONDDIR@/vz is cleared when the
# service is stopped. Therefore any edits made directly to
# @DST_CRONDDIR@/vz will be lost anytime the vz service
# restarts.
#
# To make changes, edit the master copy @CRONDDIR@/vz
# (or add your own @CRONDDIR@/vz* files)
# and run /etc/init.d/vz update-cron.
#
# The reason this is done this way is because the OpenVZ cron job
# should only be invoked if the vz service is enabled and not
# just as a consequence of installing the vzctl as it was the case
# previously. The file @DST_CRONDDIR@/vz
# cannot simply be linked to the master copy in
# @CRONDDIR@/ because for security reasons cron
# will not process crontab files that are links or writeable by
# anybody else but root, thus the file must be copied
# into @DST_CRONDDIR@/ with the right ownership and permissions.


# Clean VEs network rules (arp,route,iptables)
*/5 * * * * root @PKGDATADIR@/scripts/vpsnetclean

# Start VEs marked as rebooted.
*/5 * * * * root @PKGDATADIR@/scripts/vpsreboot
