# etc/network/if-up.d/vzifup-post.in
#!/bin/sh

[ "$ADDRFAM" = "inet" -o "$ADDRFAM" = "inet6" ] || exit 0

exec @SBINDIR@/vzifup-post "$IFACE"
