# etc/network/if-up.d/vzifup-post.in
#!/bin/sh

exec @SBINDIR@/vzifup-post "$IFACE"
